# EP1/data-workflow/crontab
*/5 * * * * /home/linuxmint/Documents/Scripts/weather_script.sh && /home/linuxmint/Documents/Scripts/stock_script.sh
